# dev-python/httpsrv/files/httpsrv
#!/sbin/openrc-run
#
# Copyright (C) 2020 Daniel Mueller <deso@posteo.net>
# SPDX-License-Identifier: GPL-3.0-or-later
#
# $Header: Serve the contents of a directory over HTTP.$

PIDFILE="/tmp/httpsrv.pid"
STDOUT="/tmp/httpsrv.log"


depend() {
  after net
}

start() {
  if [ -z "${HTTPSRV_DIRECTORY}" ]; then
    msg=$(cat <<EOF
The following variables are required for ${RC_SVCNAME} to work:
HTTPSRV_DIRECTORY: The directory the contents of which to serve.

Optionally you may want to specify:
HTTPSRV_ADDRESS: The address to listen on (e.g., '' or localhost).
HTTPSRV_USER: The user to run as (optionally provide :<group> as well).

Please define them in /etc/conf.d/${RC_SVCNAME}.
EOF
)
    eerror "${msg}"
    return 1
  fi

  start-stop-daemon --make-pidfile --pidfile "${PIDFILE}" --stdout "${STDOUT}" \
    --user "${HTTPSRV_USER:-nobody:nogroup}" --start --background --exec httpsrv -- \
    --address "${HTTPSRV_ADDRESS:-}" --directory "${HTTPSRV_DIRECTORY}"
}

stop() {
  start-stop-daemon --pidfile "${PIDFILE}" --stop
}
